# scripts/systemd/aipm-kiro-cli.service
[Unit]
Description=AIPM Kiro CLI Chat Process
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
ExecStart=/home/ec2-user/.local/bin/kiro-cli chat --trust-all-tools
Restart=always
RestartSec=10
Environment=NODE_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aipm-kiro-cli

[Install]
WantedBy=multi-user.target
