# config/aipm-backend.service
[Unit]
Description=AIPM Backend API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
EnvironmentFile=/home/ec2-user/aipm/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node apps/backend/server.js
Restart=on-failure
RestartSec=5s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
