# config/kiro-wrapper@.service
[Unit]
Description=Kiro Wrapper Session %i
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
ExecStart=/usr/bin/node scripts/kiro-wrapper.js %i 900%i
Restart=always
RestartSec=5

# Resource limits
MemoryMax=300M
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiro-wrapper-%i

[Install]
WantedBy=multi-user.target
