# config/kiro-wrapper@.service
[Unit]
Description=Kiro Wrapper Session %i
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
ExecStart=/usr/bin/node scripts/kiro-wrapper.js %i 900%i
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=300M
TasksMax=100

# Logging
StandardOutput=append:/home/ec2-user/kiro-wrapper.log
StandardError=append:/home/ec2-user/kiro-wrapper.log
SyslogIdentifier=kiro-wrapper-%i

[Install]
WantedBy=multi-user.target
