# config/kiro-session-pool.service
[Unit]
Description=Kiro CLI Session Pool Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
ExecStart=/usr/bin/node /home/ec2-user/aipm/scripts/kiro-session-pool.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/kiro-session-pool.log
StandardError=append:/var/log/kiro-session-pool-error.log

# Environment
Environment="NODE_ENV=production"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Resource limits
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
