# config/queue-cleanup-prod.service
[Unit]
Description=AIPM DynamoDB Queue Cleanup Service (Production)
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/aipm
Environment="AWS_REGION=us-east-1"
Environment="DYNAMODB_QUEUE_TABLE=aipm-semantic-api-queue-prod"
ExecStart=/usr/bin/node /home/ec2-user/aipm/scripts/queue-cleanup.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/queue-cleanup.log
StandardError=append:/var/log/queue-cleanup.log

[Install]
WantedBy=multi-user.target
