# scripts/aipm-ip-update.service
[Unit]
Description=Update EC2 IP in S3 config on boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Detect environment from instance tags
ExecStart=/bin/bash -c 'INSTANCE_ID=$(ec2-metadata --instance-id 2>/dev/null | cut -d" " -f2); ENV_TAG=$(aws ec2 describe-tags --region us-east-1 --filters "Name=resource-id,Values=$INSTANCE_ID" "Name=key,Values=Environment" --query "Tags[0].Value" --output text 2>/dev/null); if [[ "$ENV_TAG" == "production" ]]; then ENV="prod"; elif [[ "$ENV_TAG" == "development" ]]; then ENV="dev"; else echo "ERROR: Could not determine environment from tags"; exit 1; fi; /usr/local/bin/update-ec2-ip-to-s3.sh $ENV'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
